# ansible/roles/evilginx/templates/evilginx.service.j2
[Unit]
Description=Evilginx2 Phishing Framework
After=network.target

[Service]
Type=simple
WorkingDirectory=/opt/evilginx2
ExecStart=/opt/evilginx2/evilginx
Restart=always
RestartSec=10

[Install]
WantedBy=multi-user.target
